# usr/etc/profile.d/nix.sh
if [ -e '/nix/nix/nix/etc/profile.d/nix-daemon.sh' ]; then
  . '/nix/nix/nix/etc/profile.d/nix-daemon.sh'
fi
